# manage-server.service
[Unit]
Description=manage-server

[Service]
ExecStart=/bin/bash /manage/manage3server/run.sh

[Install]
WantedBy=multi-user.target
